# crontab
0/5 * * * * ubuntu cd /home/ubuntu/notification_system && pipenv run python "toad_notification_system.py" >> /var/log/cron.log 2>&1
