# crontab
0/5 * * * * python3 /notification_system/toad_notification_system.py /config/config.json