# crontab
* * * * * root python "/toad_notification_system.py" >> /var/log/cron.log 10>&1
